# tools/omx_init.in
#!/bin/sh

# chkconfig: 2345 30 70
# description: starts Open-MX driver

### BEGIN INIT INFO
# Provides:          mx
# Required-Start:    $network $syslog $remote_fs
# Required-Stop:     $network $syslog $remote_fs
# Default-Start:     3 5
# Default-Stop:      0 1 2 6
# Description:       Starts the Open-MX driver
### END INIT INFO

# Note that this is replaced by make install, not configure!
OMX_PREFIX=@omx_install_prefix@
MAIN_DEVICE_NAME=@omx_main_device_name@
MAIN_DEVICE_GROUP=@omx_main_device_group@
MAIN_DEVICE_MODE=@omx_main_device_mode@
RAW_DEVICE_NAME=@omx_raw_device_name@
RAW_DEVICE_GROUP=@omx_raw_device_group@
RAW_DEVICE_MODE=@omx_raw_device_mode@
OMX_STATIC_MAPPING=@omx_static_peers@
OMX_PEERS_FILE=@omx_peers_file@

OMX_CONFIG=/etc/open-mx/open-mx.conf
LINUX_RELEASE=$(uname -r)

set -e

if test ! -d ${OMX_PREFIX}; then
    echo "Something bad happened with the install script"
    echo "OMX_PREFIX isn't pointing to a valid directory"
    exit 1
fi

if test -f ${OMX_CONFIG}; then
    . ${OMX_CONFIG}
fi

usage () {
	echo "Usage: $0 [options] <command> [module load params]"
	echo " Commands:	start stop status restart"
	echo " Options:"
	echo "  --dynamic-peers		Force dynamic peer discovery"
	echo "  --static-peers		Force static peers from $OMX_PEERS_FILE"
	echo "  --static-peers=<file>		Force static peers from <file>"
	exit 1
}

[ $# -eq 0 ] && usage

while [ ${1##--} != $1 ] ; do
    if [ $1 = "--dynamic-peers" ] ; then
	OMX_STATIC_MAPPING=0
    else if [ $1 = "--static-peers" ] ; then
	OMX_STATIC_MAPPING=1
    else if [ ${1##--static-peers=} != $1 ] ; then
	OMX_STATIC_MAPPING=1
	OMX_PEERS_FILE=${1##--static-peers=}
    else
	usage
    fi fi fi
    shift
    [ $# -eq 0 ] && usage
done

OMX_MODULE_PARAMS="${OMX_MODULE_PARAMS} ifnames=${IFACES}"

# Remove the Open-MX module if loaded
unload_omx() {
    if grep open[-_]mx /proc/modules > /dev/null; then
        echo "Removing Open-MX driver"
        /sbin/rmmod open-mx
    fi
    rm -f "$MAIN_DEVICE_NAME"
}

# Load the Open-MX module
load_omx() {
    echo "Loading Open-MX driver"
    params=" $OMX_MODULE_PARAMS $*"
    module="${OMX_PREFIX}/lib/modules/${LINUX_RELEASE}/open-mx.ko"

    # modprobe the dependencies
    depends=$(modinfo $module | grep ^depends | awk '{print $2}' | sed 's/,/ /')
    for depend in $depends ; do
	echo "Loading dependency module: $depend"
	modprobe $depend
    done

    # load the driver
    /sbin/insmod $module $params

    # sleep a bit so that udev is done
    sleep 1

    # create the main device node if necessary
    if ! test -e $MAIN_DEVICE_NAME; then
	echo "Creating '$MAIN_DEVICE_NAME'"
	major=$(sed -nr -e '/^Character devices:/,/^Block devices:/s/ *([0-9]+) *misc/\1/p' /proc/devices)
	minor=$(sed -nr -e 's/ *([0-9]+) *open-mx$/\1/p' /proc/misc)
	if test -z "$major" -o -z "$minor"; then
	    echo "[ERROR] Failed to retrieve main device major/minor numbers"
	else
	    mknod "$MAIN_DEVICE_NAME" c "$major" "$minor" \
		|| echo "[ERROR] Failed to create main device node"
	fi
    fi

    # set main device node mode and group
    echo "Setting group '$MAIN_DEVICE_GROUP' and mode '$MAIN_DEVICE_MODE' on main device '$MAIN_DEVICE_NAME'"
    chgrp "$MAIN_DEVICE_GROUP" "$MAIN_DEVICE_NAME" \
	|| echo "[ERROR] Failed to set main device mode"
    chmod "$MAIN_DEVICE_MODE" "$MAIN_DEVICE_NAME" \
	|| echo "[ERROR] Failed to set main device group"

    # create the raw device node if necessary
    if ! test -e $RAW_DEVICE_NAME; then
	echo "Creating '$RAW_DEVICE_NAME'"
	major=$(sed -nr -e '/^Character devices:/,/^Block devices:/s/ *([0-9]+) *misc/\1/p' /proc/devices)
	minor=$(sed -nr -e 's/ *([0-9]+) *open-mx-raw$/\1/p' /proc/misc)
	if test -z "$major" -o -z "$minor"; then
	    echo "[ERROR] Failed to retrieve raw device major/minor numbers"
	else
	    mknod "$RAW_DEVICE_NAME" c "$major" "$minor" \
		|| echo "[ERROR] Failed to create raw device node"
	fi
    fi

    # set raw device node mode and group
    echo "Setting group '$RAW_DEVICE_GROUP' and mode '$RAW_DEVICE_MODE' on raw device '$RAW_DEVICE_NAME'"
    chgrp "$RAW_DEVICE_GROUP" "$RAW_DEVICE_NAME" \
	|| echo "[ERROR] Failed to set raw device mode"
    chmod "$RAW_DEVICE_MODE" "$RAW_DEVICE_NAME" \
	|| echo "[ERROR] Failed to set raw device group"
}

# Start the mapper
start_mapper() {
    if test ${OMX_STATIC_MAPPING} -eq 1 ; then
	# initialize peers
	if test -f ${OMX_PEERS_FILE}; then
	    echo "Initializing peers using file '${OMX_PEERS_FILE}'"
	    ${OMX_PREFIX}/bin/omx_init_peers -c "${OMX_PEERS_FILE}"
	else
	    echo "Peers file ${OMX_PEERS_FILE} does not exist, remember to run omx_peers_init with the correct file"
	fi
    else
	echo "Starting the dynamic peer mapper"
	${OMX_PREFIX}/bin/omxoed &
	sleep 2
    fi
}

# Stop the mapper
stop_mapper() {
    if grep open[-_]mx /proc/modules > /dev/null; then
	echo "Killing the dynamic peer mapper"
	killall -9 omxoed || /bin/true
    fi
}

cmd="$1"
shift

case "$cmd" in
    start)
	load_omx "$@"
	start_mapper
	;;
    stop)
	stop_mapper
	unload_omx
	;;
    restart)
	stop_mapper
	unload_omx
	load_omx "$@"
	start_mapper
        ;;
    start-mapper)
	start_mapper
	;;
    stop-mapper)
	stop_mapper
	;;
    restart-mapper)
	stop_mapper
	start_mapper
	;;
    status)
	if grep open[-_]mx /proc/modules > /dev/null; then
	    echo "Open-MX driver is loaded"
	else
	    echo "Open-MX driver is not loaded"
	fi
	;;
    *)
	usage
esac

exit 0
